# root/etc/s6-overlay/s6-rc.d/prowlarr/run
#!/usr/bin/execlineb -P

/opt/prowlarr/bin/Prowlarr -nobrowser -data=/container/data
